# service/dds661-polling.service
[Unit]
Description=DDS661 Polling Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/dds661
ExecStart=/opt/dds661/.venv/bin/python3 /opt/dds661/polling.py --config /opt/dds661/config.yaml
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
